# Masonry-Code/PythonLibImport/miniforge/miniforge3/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/alisaleh/miniforge3/bin/conda";
setenv _CONDA_EXE "/Users/alisaleh/miniforge3/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/alisaleh/miniforge3/bin/python";
setenv _CONDA_ROOT "/Users/alisaleh/miniforge3";

source "/Users/alisaleh/miniforge3/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"